# Zoom/change_license_type/example_cron_structure/etc/cron.d/change_license_cron
# Run job at 05:00 every day
0 5 * * *   root   /local/zoom/scripts/change_license.sh
